# docker/light/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:redis]
command=redis-server --appendonly yes --dir /build/data
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=node /build/apps/api/dist/main.js
directory=/build/apps/api
user=mitshe
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=10
environment=NODE_ENV="production",PORT="3001",DATABASE_URL="file:/build/data/mitshe.db",REDIS_URL="redis://localhost:6379"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=npm start
directory=/build/apps/web
user=mitshe
autostart=true
autorestart=true
priority=30
startretries=5
startsecs=10
environment=NODE_ENV="production",PORT="3000",NEXT_PUBLIC_API_URL="http://localhost:3001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
